# templates/supervisor.conf.j2
[program:helga]
command=/opt/helga/bin/helga
stopsignal=INT
stopwaitsecs=16
environment=HELGA_SETTINGS=custom_settings
